# docker/s6/s6-rc.d/web/run
#!/command/with-contenv sh
exec frankenphp php-server --listen :8080 --root /app/public
