# docker/s6/s6-rc.d/web/run
#!/command/execlineb -P
frankenphp php-server --listen :8080 --root /app/public
