# docker/s6/s6-rc.d/init-perms/up
/command/with-contenv sh -c "mkdir -p /data/covers /config && chown -R www-data:www-data /data /config"
