# docker/s6/s6-rc.d/worker4/run
#!/command/with-contenv sh
# Extra queue worker вЂ” slot 4 of 4. Runs a queue:work only when QUEUE_WORKERS >= 4
# (default 1); otherwise idles so s6 keeps supervising it without restart-thrashing.
# Slot 1 is the always-on `worker` service. / дЅ™е‰°гѓЇгѓјг‚«гѓј(slot4)гЂ‚QUEUE_WORKERS>=4гЃ§иµ·е‹•гЂЃжњЄжєЂгЃЇеѕ…ж©џгЂ‚
if [ "${QUEUE_WORKERS:-1}" -ge 4 ] 2>/dev/null; then
    exec s6-setuidgid www-data php /app/artisan queue:work --sleep=3 --tries=3 --max-time=3600
fi
exec sleep infinity
