# docker/s6/s6-rc.d/scheduler/run
#!/command/with-contenv sh
exec s6-setuidgid www-data php /app/artisan schedule:work
