# docker/s6/s6-rc.d/scheduler/run
#!/command/with-contenv sh
exec php /app/artisan schedule:work
