# docker/s6/s6-rc.d/worker/run
#!/command/with-contenv sh
exec s6-setuidgid www-data php /app/artisan queue:work --sleep=3 --tries=3 --max-time=3600
